# wa-bot.service
[Unit]
Description=WhatsApp PC Control Bot
After=network.target

[Service]
Type=simple
User=kristoffersanio
WorkingDirectory=/Users/kristoffersanio/git/jkristoffer.com/wa-pc-bot
ExecStart=/usr/bin/node index.js
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
